# roles/ghost_server/templates/ghost.service.j2
## ghost systemd unit
[Unit]
Description=Ghost Service {{ item.hostname }}
Requires=docker.service
After=docker.service
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
# give docker time to pull the image if it's not already available
# likely not necessary once we get past the PVE7 docker on LXC on ZFS bug
# failed to register layer: unlinkat /tmp/v8-compile-cache-0/10.2.154.26-node.28: invalid argument
# TimeoutStartSec=300
# ExecStartPre=/usr/bin/docker pull {{ ghost_image }}
User=ghost
Restart=always
ExecStart=/usr/bin/docker run \
  --rm \
  --name ghost-{{ item.hostname }} \
  --network ghost \
  --user {{ ghost_uid }}:{{ ghost_gid }} \
  -e NODE_ENV=production \
  -e url="https://{{ item.hostname }}" \
  -e mail__transport=SMTP \
  -e mail__options__service="smarthost" \
  -e mail__options__host="{{ smarthost_ip }}" \
  -e mail__options__port=25 \
{% if item.smtp_from is defined %}
  -e mail__from="{{ item.smtp_from }}" \
{% endif %}
  -e database__client=sqlite3 \
  -e database__connection__filename="content/data/ghost.db" \
  -e database__useNullAsDefault=true \
  -e database__debug=false \
  -p {{ item.port }}:2368 \
  -v /home/ghost/{{ item.hostname }}:/var/lib/ghost/content \
  {{ ghost_image }}
ExecStop=/usr/bin/docker stop ghost-{{ item.hostname }}

[Install]
WantedBy=default.target

